# auto-shutdown.service
[Unit]
Description=Proxmox Auto-Shutdown Monitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/serverscripts
ExecStart=/usr/bin/python3 "/root/serverscripts/auto stop.py"
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
